# icecast.service
[Unit]
Description=Icecast2

[Service]
Type=simple
ExecStart=/usr/bin/icecast.sh start
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
